# gentoo/files/bcfg2-client
#!/sbin/runscript
#
# bcfg - bcfg configuration client
#
# 

depend () {
	need net
}

start () {
	ebegin "Starting bcfg"
	if [ -e /etc/bcfg.dont_run ]; then
		einfo "bcfg does not need to run"
		eend 0
	else
		start-stop-daemon --start --quiet \
			--startas /usr/sbin/bcfg.py -- -q -s
		eend $? "Failed to start bcfg"
	fi
}
stop () {
	ebegin "Stopping bcfg"
	/bin/true
	eend 0
}
