# debian/bcfg2.cron.daily
#!/bin/sh
BCFG2CRON=
if [ -x /usr/libexec/bcfg2-cron ]; then
    BCFG2CRON=/usr/libexec/bcfg2-cron
elif [ -x /usr/lib/bcfg2/bcfg2-cron ]; then
    BCFG2CRON=/usr/lib/bcfg2/bcfg2-cron
elif type bcfg2-cron 2>&1 >/dev/null; then
    BCFG2CRON=bcfg2-cron
else
    echo "No bcfg2-cron command found"
    exit 1
fi
$BCFG2CRON --daily 2>&1 | logger -t bcfg2-cron -p daemon.info -i
