# debian/bcfg2.cron.hourly
#!/bin/sh
[ -x /usr/lib/bcfg2/bcfg2-cron ] || exit 0
/usr/lib/bcfg2/bcfg2-cron --hourly > /dev/null 2>&1 || true
